# etc/autostart/launchy.desktop
[Desktop Entry]
Type=Application
Terminal=false
Exec=/usr/bin/launchy
Name=launchy

